# src/rootfs/etc/services.d/unifi/run
#!/usr/bin/sh

cd /var/lib/unifi
exec /usr/bin/java -Xmx1024M -jar /usr/lib/unifi/lib/ace.jar start
